# contrib/systemd/smolotchi-display.service
[Unit]
Description=Smolotchi Display (Waveshare e-paper)
After=network.target

[Service]
Type=simple
User=smolotchi
WorkingDirectory=/opt/smolotchi
Environment=PYTHONUNBUFFERED=1
Environment=SMOLOTCHI_DEVICE=pi_zero
Environment=SMOLOTCHI_DB=/var/lib/smolotchi/events.db
Environment=SMOLOTCHI_ARTIFACT_ROOT=/var/lib/smolotchi/artifacts
ExecStart=/opt/smolotchi/.venv/bin/python -m smolotchi.cli display --device pi_zero --db /var/lib/smolotchi/events.db --artifact-root /var/lib/smolotchi/artifacts
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
