# packaging/systemd/smolotchi-ai-worker.service
[Unit]
Description=Smolotchi AI Worker
After=network.target

[Service]
User=pi
Environment=PYTHONUNBUFFERED=1
Environment=SMOLOTCHI_PROJECT_DIR=/home/pi/smolotchi
Environment=SMOLOTCHI_VENV_PYTHON=/home/pi/smolotchi/.venv/bin/python
EnvironmentFile=-/etc/smolotchi/smolotchi.env
WorkingDirectory=${SMOLOTCHI_PROJECT_DIR}
ExecStart=${SMOLOTCHI_VENV_PYTHON} -m smolotchi.ai.worker --loop --log-level INFO
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
