# packaging/systemd/smolotchi-display.service
[Unit]
Description=Smolotchi E-Paper Display
After=smolotchi-core.service
Requires=smolotchi-core.service

[Service]
Type=simple
User=smolotchi
Group=smolotchi
WorkingDirectory=/home/smolotchi
EnvironmentFile=-/etc/smolotchi/env
ExecStart=/usr/local/bin/smolotchi display
Restart=always
RestartSec=15

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
