# packaging/systemd/smolotchi-core.service
[Unit]
Description=Smolotchi Core (WiFi / LAN / Tools Engine)
After=network.target
Wants=network.target
Conflicts=smolotchi-core-net.service

[Service]
Type=notify
User=smolotchi
Group=smolotchi
WorkingDirectory=/var/lib/smolotchi
EnvironmentFile=-/etc/smolotchi/env
ExecStart=/usr/local/bin/smolotchi core
Restart=always
RestartSec=5
WatchdogSec=30
NotifyAccess=main

MemoryMax=350M
CPUQuota=70%

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
