# packaging/systemd/smolotchi-prune.service
[Service]
Type=oneshot
User=smolotchi
Group=smolotchi
EnvironmentFile=-/etc/smolotchi/env
ExecStart=/usr/local/bin/smolotchi prune
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/smolotchi /run/smolotchi
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true
