# packaging/systemd/smolotchi-prune.service
[Unit]
Description=Smolotchi Artifact Retention Prune

[Service]
Type=oneshot
User=smolotchi
Group=smolotchi
EnvironmentFile=/etc/smolotchi/env
ExecStart=/usr/local/bin/smolotchi prune

# Optional hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/var/lib/smolotchi

[Install]
WantedBy=multi-user.target
