# packaging/systemd/smolotchi-prune.service
[Service]
Type=oneshot
ExecStart=/home/smolotchi/.venv/bin/python -m smolotchi.cli prune
