# packaging/systemd/smolotchi-core.service
[Unit]
Description=Smolotchi Core (WiFi / LAN / Tools Engine)
After=network.target
Wants=network.target

[Service]
Type=notify
User=smolotchi
Group=smolotchi
WorkingDirectory=/var/lib/smolotchi
EnvironmentFile=-/etc/smolotchi/env
ExecStart=/usr/local/bin/smolotchi core
Restart=always
RestartSec=5
WatchdogSec=30
NotifyAccess=main

AmbientCapabilities=CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_RAW

NoNewPrivileges=true

PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/smolotchi /run/smolotchi
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true

MemoryMax=350M
CPUQuota=70%

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
