# packaging/systemd/smolotchi-core.service
[Unit]
Description=Smolotchi Core Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=smolotchi
Group=smolotchi
WorkingDirectory=/opt/smolotchi
EnvironmentFile=/etc/smolotchi/smolotchi.env
RuntimeDirectory=smolotchi
RuntimeDirectoryMode=0775
ExecStartPre=/bin/mkdir -p /run/smolotchi/locks
ExecStart=/opt/smolotchi/.venv/bin/python -m smolotchi.cli --config /etc/smolotchi/config.toml core
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
