# packaging/systemd/smolotchi-web.service
[Unit]
Description=Smolotchi Web UI
After=network.target smolotchi-core.service
Requires=smolotchi-core.service

[Service]
Type=simple
User=smolotchi
Group=smolotchi
WorkingDirectory=/home/%u/smolotchi
EnvironmentFile=/etc/smolotchi/env
ExecStart=/home/%u/smolotchi/.venv/bin/python -m smolotchi.cli web
Restart=always
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/var/lib/smolotchi

[Install]
WantedBy=multi-user.target
