# packaging/systemd/smolotchi-prune.timer
[Unit]
Description=Smolotchi Artifact Retention

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
