# packaging/systemd/smolotchi-ai.service
[Unit]
Description=Smolotchi AI Worker
After=smolotchi-core.service
Requires=smolotchi-core.service

[Service]
Type=simple
User=smolotchi
Group=smolotchi
WorkingDirectory=/var/lib/smolotchi
EnvironmentFile=-/etc/smolotchi/env
RuntimeDirectory=smolotchi smolotchi/locks
RuntimeDirectoryMode=0775
StateDirectory=smolotchi
StateDirectoryMode=0775
ExecStart=/usr/local/bin/smolotchi ai
Restart=always

[Install]
WantedBy=multi-user.target
