# packaging/systemd/smolotchi-web.service
[Unit]
Description=Smolotchi Web UI
After=network-online.target smolotchi-core.service
Wants=network-online.target
Requires=smolotchi-core.service

[Service]
Type=simple
User=smolotchi
Group=smolotchi

WorkingDirectory=%h/smolotchi
EnvironmentFile=/etc/smolotchi/env

ExecStart=/usr/local/bin/smolotchi web

Restart=always
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/var/lib/smolotchi

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
