# packaging/systemd/smolotchi-ai.service
[Unit]
Description=Smolotchi AI Worker
After=smolotchi-core.service
Requires=smolotchi-core.service

[Service]
Type=simple
User=smolotchi
Group=smolotchi
WorkingDirectory=/home/smolotchi
EnvironmentFile=/etc/smolotchi/env
ExecStart=/home/smolotchi/.venv/bin/python -m smolotchi.ai.worker --loop --log-level INFO
Restart=always
RestartSec=10

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/smolotchi

[Install]
WantedBy=multi-user.target
