# packaging/systemd/smolotchi-display.service
[Unit]
Description=Smolotchi E-Paper Display
After=smolotchi-core.service
Requires=smolotchi-core.service

[Service]
Type=simple
User=smolotchi
Group=smolotchi
WorkingDirectory=/home/smolotchi
EnvironmentFile=-/etc/smolotchi/env
RuntimeDirectory=smolotchi smolotchi/locks
RuntimeDirectoryMode=0775
StateDirectory=smolotchi
StateDirectoryMode=0775
ExecStart=/usr/local/bin/smolotchi display
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target
