# packaging/systemd/smolotchi-web.service
[Unit]
Description=Smolotchi Web UI
After=network.target smolotchi-core.service
Requires=smolotchi-core.service

[Service]
Type=simple
User=smolotchi
Group=smolotchi
WorkingDirectory=/var/lib/smolotchi
EnvironmentFile=-/etc/smolotchi/env
ExecStart=/usr/local/bin/smolotchi web
Restart=always

[Install]
WantedBy=multi-user.target
