# packaging/systemd/smolotchi-ai.service
[Unit]
Description=Smolotchi AI Worker
After=network-online.target smolotchi-core.service
Wants=network-online.target
Requires=smolotchi-core.service

[Service]
Type=simple
User=smolotchi
Group=smolotchi

EnvironmentFile=/etc/smolotchi/env

ExecStart=/usr/local/bin/smolotchi ai --loop --log-level INFO

Restart=always
RestartSec=10

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/var/lib/smolotchi

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
