# packaging/systemd/smolotchi-display.service
[Unit]
Description=Smolotchi E-Paper Display
After=network-online.target smolotchi-core.service
Wants=network-online.target
Requires=smolotchi-core.service

[Service]
Type=simple
User=smolotchi
Group=smolotchi
EnvironmentFile=/etc/smolotchi/env

RuntimeDirectory=smolotchi smolotchi/locks
RuntimeDirectoryMode=0775
StateDirectory=smolotchi
StateDirectoryMode=0775

ExecStart=/usr/local/bin/smolotchi display

Restart=always
RestartSec=15

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/var/lib/smolotchi /run/smolotchi

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
