# packaging/systemd/dropins/smolotchi-display.service.d/20-service-paths.conf
[Service]
ProtectHome=false
PrivateNetwork=false
ReadWritePaths=/var/lib/smolotchi /run/smolotchi
