# packaging/systemd/dropins/smolotchi-display.service.d/20-service-paths.conf
[Service]
ReadWritePaths=/var/lib/smolotchi /run/smolotchi
# Display runs without ambient capabilities.
AmbientCapabilities=
CapabilityBoundingSet=
