# packaging/systemd/dropins/smolotchi-prune.service.d/20-service-paths.conf
[Service]
ReadWritePaths=/var/lib/smolotchi /run/smolotchi
