# packaging/systemd/dropins/smolotchi-prune.service.d/20-service-paths.conf
[Service]
ProtectHome=true
PrivateNetwork=true
ReadWritePaths=/var/lib/smolotchi /run/smolotchi
