# packaging/systemd/dropins/smolotchi-ai.service.d/20-service-paths.conf
[Service]
ProtectHome=true
ReadWritePaths=/var/lib/smolotchi /run/smolotchi
