# packaging/systemd/dropins/smolotchi-web.service.d/20-service-paths.conf
[Service]
PrivateNetwork=false
ReadWritePaths=/var/lib/smolotchi /run/smolotchi
