# packaging/systemd/dropins/smolotchi-web.service.d/20-service-paths.conf
[Service]
ProtectHome=true
ReadWritePaths=/var/lib/smolotchi /run/smolotchi
# Web UI runs without ambient capabilities.
AmbientCapabilities=
CapabilityBoundingSet=
